# resources/autostart/pinsnap.desktop
[Desktop Entry]
Type=Application
Name=PinSnap
Comment=A PixPin-like screenshot, annotation, and pin tool for Linux
Exec=pinsnap
Icon=pinsnap
Terminal=false
Categories=Utility;Graphics;
Keywords=screenshot;capture;annotate;pin;screen;
StartupNotify=false
X-GNOME-Autostart-enabled=true
X-MATE-Autostart-enabled=true